# conf/inittab
n1:2345:respawn:/usr/local/ns/bin/nsd.start

